# roles/cve_2026_31431/templates/ebpf_loader.service.j2
[Unit]
Description=CVE-2026-31431 eBPF LSM Protection
Documentation=https://github.com/parmstro/cfDr
After=network.target
Before=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/bpftool prog load {{ ebpf_install_dir }}/{{ ebpf_program_name }}.bpf.o {{ bpf_fs }}/{{ ebpf_program_name }} type lsm
ExecStop=/bin/rm -f {{ bpf_fs }}/{{ ebpf_program_name }}

# Security hardening
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ bpf_fs }}
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
